# systemd/agent-terminal-starthub.service
[Unit]
Description=Agent terminal nudge for starthub tmux session

[Service]
Type=oneshot
ExecStart=/home/mikhail/Desktop/shadow/scripts/tmux-agent-work-ping starthub
