# systemd/agent-terminal@.service.template
[Unit]
Description=Agent terminal dispatch for %i OpenCode session

[Service]
Type=oneshot
EnvironmentFile=/home/mikhail/.config/shadow/opencode.env
ExecStart=/home/mikhail/Desktop/shadow/scripts/opencodectl agent-ping %i
