# systemd/agent-terminal-nixelo.service
[Unit]
Description=Agent terminal nudge for nixelo tmux session

[Service]
Type=oneshot
ExecStart=/home/mikhail/Desktop/shadow/scripts/tmux-agent-work-ping nixelo
