# systemd/manual-terminal@.service.template
[Unit]
Description=Manual OpenCode dispatch for %i repo

[Service]
Type=oneshot
ExecStart=/home/mikhail/Desktop/shadow/scripts/opencodectl manual-ping %i
