# systemd/manual-terminal-nixelo.service
[Unit]
Description=Manual OpenCode dispatch for nixelo repo

[Service]
Type=oneshot
ExecStart=/home/mikhail/Desktop/shadow/scripts/opencodectl manual-ping nixelo
