# systemd/manual-terminal-nixelo.service
[Unit]
Description=Manual terminal nudge for nixelo tmux session

[Service]
Type=oneshot
ExecStart=/home/mikhail/Desktop/shadow/scripts/tmux-manual-work-ping nixelo
