# systemd/manual-terminal-nixelo.service
[Unit]
Description=Manual terminal nudge for nixelo tmux session

[Service]
Type=oneshot
EnvironmentFile=/home/mikhail/.config/shadow/opencode.env
ExecStart=/home/mikhail/Desktop/shadow/scripts/opencodectl manual-ping nixelo
