# systemd/opencode-auto-nixelo.timer
[Unit]
Description=Run OpenCode auto-cycle for nixelo every minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
Unit=opencode-auto-nixelo.service

[Install]
WantedBy=timers.target
