# systemd/prci-terminal-nixelo.timer
[Unit]
Description=Run PR-CI terminal nudge for nixelo every 1 minute

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
Unit=prci-terminal-nixelo.service

[Install]
WantedBy=timers.target