# systemd/prci-terminal-nixelo.service
[Unit]
Description=PR-CI terminal nudge for nixelo tmux session

[Service]
Type=oneshot
ExecStart=/home/mikhail/Desktop/shadow/scripts/opencodectl prci-ping nixelo
