# systemd/prci-terminal@.service.template
[Unit]
Description=PR-CI terminal nudge for %i tmux session

[Service]
Type=oneshot
ExecStart=/home/mikhail/Desktop/shadow/scripts/opencodectl prci-ping %i
