# systemd/watcher.service
[Unit]
Description=Terminal automation watcher - deterministic health checks
After=network.target

[Service]
Type=oneshot
ExecStart=/home/mikhail/Desktop/shadow/scripts/watcher.sh
TimeoutStartSec=60
